# dotut.service
[Unit]
Description = Digital Ocean Web App service

[Service]
ExecStart=/home/goapp/src/rest/app.sh

[Install]
WantedBy = multi-user.target
